# cmdeploy/src/cmdeploy/metrics.cron.j2
*/5 * * * * root {{ config.execpath }} {{ config.mailboxes_dir }} >/var/www/html/metrics
