# cmdeploy/src/cmdeploy/metrics.cron.j2
*/5 * * * * root {{ config.execpath }} /home/vmail/mail/{{ config.mail_domain }} >/var/www/html/metrics
