# cmdeploy/src/cmdeploy/service/expunge.timer
[Unit]
Description=Run expunge.service daily

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
