# cmdeploy/src/cmdeploy/service/metrics.service.j2
[Unit]
Description=Generate metrics in /var/www/html/metrics

[Service]
ExecStart={{ config.execpath }} /home/vmail/mail/{{ config.mail_domain }} > /var/www/html/metrics
