# cmdeploy/src/cmdeploy/service/metrics.service.j2
[Unit]
Description=Generate metrics in /var/www/html/metrics every 5 minutes

[Service]
WorkingDirectory=/var/www/magento-website/current
ExecStart=/usr/bin/php bin/magento indexer:reindex
